# internet-canary.conf

start on runlevel [2345]

script
    exec /usr/bin/internet-canary.py &>> /tmp/internet-canary.log
end script
