# internet-canary.conf

start on runlevel [2345]

script
	exec /usr/bin/internet-canary.py /etc/internet-canary.d/internet-canary.json
endscript
